# src/CashVault/deb-package/DEBIAN/prerm
#!/bin/bash
systemctl stop monroe-backbone
systemctl disable monroe-backbone
